# to-symlink/podman-containers/btrfs-backup.timer
[Unit]
Description=Run btrfs backup daily

[Timer]
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
