# to-symlink/podman-containers/config-backup.service
[Unit]
Description=Config backup script

[Service]
Type= oneshot
ExecStart= /home/sayantan/home-manager-nix/scripts/config-backups.sh
